# cron/crontab
0 0 1 * * zpool scrub bucket
0 4 * * * $HOME/src/gh-backup/backup.py >> $HOME/src/gh-backup/log 2>&1